# res/out/ubuntu/etc/systemd/system/plug2proxy.service
[Unit]
Description=Plug2Proxy
After=network.target

[Service]
User=root
ExecStart=/usr/sbin/plug2proxy /etc/plug2proxy/config.json --data-dir /etc/plug2proxy
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
